# dyndns.timer
[Unit]
Description=Update Route53 DNS records

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
